# templates/etc/network/if-up.d/ip6tables_load.j2
#!/bin/bash

# This file was generated by Ae for {{ansible_fqdn}}
# Do NOT modify this file by hand!

{{ip6tables_rules_path}}
exit $?
